# ppp.conf
[program:ppp]
command=python3 /home/curator2/ppp/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/ppp/test.err.log
stdout_logfile=/home/curator2/ppp/test.out.log